# mod2/systemd/comment.service
[Unit]
Description=mod2 comment bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/trilobyte/ai/mod2
EnvironmentFile=/home/trilobyte/ai/mod2/.env
ExecStart=/usr/bin/env python3 comment.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
