# bibleread/systemd/pipeline.service
[Unit]
Description=bibleread pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/trilobyte/ai/bibleread
EnvironmentFile=/home/trilobyte/ai/bibleread/.env
ExecStart=/usr/bin/env python3 app.py
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
