# bibleread/systemd/pipeline.timer
[Unit]
Description=bibleread pipeline schedule

[Timer]
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
